# deploy/orquester.service
[Unit]
Description=Orquester daemon
After=network-online.target
Wants=network-online.target

[Service]
User=orquester
Group=orquester
WorkingDirectory=/opt/orquester
EnvironmentFile=/etc/orquester/daemon.env
# tsx writes its transpile cache to TMPDIR; ProtectSystem=strict makes /tmp
# unavailable, so point it at a writable, persistent dir under the appdir.
Environment=TMPDIR=/var/lib/orquester/tmp
ExecStartPre=/usr/bin/mkdir -p /var/lib/orquester/tmp
# Agent installs ("npm install -g ...") must land somewhere writable. The default
# global prefix is /usr, which is read-only under ProtectSystem=strict and not
# writable by the non-root service user, so installs fail with "npm error code
# ENOENT". Point npm at a prefix inside ReadWritePaths and put its bin/ on PATH
# so the daemon can also resolve and launch the agents it installs.
Environment=NPM_CONFIG_PREFIX=/var/lib/orquester/.npm-global
# ~/.local/bin first: Claude Code's native installer (install.sh / `claude update`) manages a
# symlink there into ~/.local/share/claude/versions/<v>, flipped atomically on update — so a
# running agent's binary is never rewritten in place. (npm -g rewrites the binary in place, which
# races the agent's own "binary changed → self-restart" and dies with EACCES/ENOENT mid-swap.)
# .npm-global/bin still follows for the other npm-installed agents (codex, opencode, ...).
Environment=PATH=/var/lib/orquester/.local/bin:/var/lib/orquester/.npm-global/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStartPre=/usr/bin/mkdir -p /var/lib/orquester/.npm-global/bin
# This monorepo is noEmit (TS is run via tsx, same as `pnpm dev:daemon`); there is
# no standalone dist build, so run the daemon through tsx as a single node process.
ExecStart=/usr/bin/node --import tsx /opt/orquester/apps/daemon/src/cli.ts --appdir /var/lib/orquester
Restart=always
RestartSec=2
# CRITICAL (Phase 2): only signal the main node process on stop/restart, so the
# daemonized tmux server (and its panes) survive a daemon restart.
KillMode=process
# Sessions run package managers via scoped sudo (/etc/sudoers.d/orquester-pkg, installed by
# deploy/provision-devtools.sh). NoNewPrivileges would block setuid elevation for the daemon's
# whole process subtree, so sudo cannot run at all while it is on — it must be off.
NoNewPrivileges=false
# Loosened from bare `strict` so scoped `sudo apt` can write /usr,/etc,/var. Kept as `strict` WITH
# explicit carve-outs (not `off`) so /opt (daemon code), /boot, /root, /home stay read-only even to
# a root session. DAC still gates writes: the unprivileged daemon user can't touch these paths; only
# root-via-sudo-apt can. If a maintainer script ever hits a read-only path, switch this to `off`.
ProtectSystem=strict
ReadWritePaths=/var/lib/orquester /usr /etc /var /run /tmp
ProtectHome=true
# tmux socket lives under /var/lib/orquester, not /tmp:
PrivateTmp=false

[Install]
WantedBy=multi-user.target
